# docs/cf/ddns.service
[Service]
ExecStart = %h/.local/bin/ddns dns update
TimeoutSec = 1min

[Unit]
Description = DDNS
