# docs/ddns/ddns.service
[Service]
ExecStart = %h/.local/bin/ddns update
TimeoutSec = 1min

[Unit]
Description = DDNS
